# aleo.service
[Unit]
Description="Aleo systemd service"
After=network.target
[Service]
User=aleo_service
Restart=always
RestartSec=10
EnvironmentFile=/etc/aleo/aleo-service.env
ExecStart=/var/lib/aleo-data/snarkos --connect 138.201.190.40:4131,135.181.200.109:4131,34.68.94.165:4131,188.40.185.61:4131,49.12.66.192:4131,157.90.22.38:4131,34.71.37.215:4131,135.180.128.248:4131,161.35.22.1:4131,95.216.208.150:4131,94.130.171.118:4131,95.217.9.46:4131,78.46.122.47:4131,95.58.110.64:4131,176.49.150.214:4131,192.248.186.23:4131,49.12.109.19:4131,34.89.53.207:4131,77.232.40.242:4131,91.122.52.83:4131,65.21.2.130:4131,95.217.23.84:4131,159.69.119.239:4131,134.209.75.17:4131,34.101.177.31:4131,116.203.99.87:4131,34.70.100.68:4131,157.90.232.92:4131,188.40.185.58:4131,95.179.254.30:4131,85.114.142.245:4131,135.181.144.152:4131,167.172.34.225:4131,91.211.249.166:4131,135.181.152.251:4131,178.166.225.186:4131,95.217.163.88:4131,77.222.52.121:4131,108.61.247.4:4131,85.14.222.140:4131,116.202.13.2:4131
[Install]
WantedBy=multi-user.target
